# services/system/beo-spotify-fetch.timer
[Unit]
Description=Hourly Spotify playlist fetch for BeoSound5c

[Timer]
OnCalendar=*:00:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
