# services/system/beo-input.service
[Unit]
Description=BeoSound5c Input Hardware Server for USB events
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/kirsten/beosound5c/services
ExecStart=/usr/bin/python3 /home/kirsten/beosound5c/services/input.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 