# services/system/beo-bluetooth.service
[Unit]
Description=BeoSound5c Bluetooth Remote service (HID)
After=network.target bluetooth.target
Wants=network.target bluetooth.target
Requires=bluetooth.target

[Service]
Type=simple
# Requires user to be in the 'input' group for /dev/input access.
# install.sh adds kirsten to the input group automatically.
User=kirsten
SupplementaryGroups=input
WorkingDirectory=/home/kirsten/beosound5c/services
EnvironmentFile=-/etc/beosound5c/config.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/kirsten/beosound5c/services/bluetooth.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
