# services/system/beo-masterlink.service
[Unit]
Description=BeoSound5c MasterLink Sniffer + Mixer HTTP API
After=network.target
Wants=network.target

[Service]
Type=simple
User=kirsten
Group=kirsten
WorkingDirectory=/home/kirsten/beosound5c/services
EnvironmentFile=-/etc/beosound5c/config.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/kirsten/beosound5c/services/masterlink.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 
