# services/system/beo-sonos.service
[Unit]
Description=BeoSound5c Sonos Player (artwork, metadata, monitoring)
After=network.target
Wants=network.target

[Service]
Type=simple
User=kirsten
Group=kirsten
WorkingDirectory=/home/kirsten/beosound5c/services
EnvironmentFile=-/etc/beosound5c/secrets.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/kirsten/beosound5c/services/players/sonos.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 
