# system/beo-sniffer.service
[Unit]
Description=BeoRemote Sniffer (sniffer.py)
After=network.target bluetooth.service

[Service]
Type=simple
User=kirsten
WorkingDirectory=/home/kirsten/beosound5c/hw
ExecStart=/usr/bin/python3 /home/kirsten/beosound5c/hw/sniffer.py
Restart=always
RestartSec=5
# Send all stdout/stderr into journald (kept in RAM)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
