# system/beo-webclient.service
[Unit]
Description=BeoSound5c WebClient (Chromium Kiosk)
After=multi-user.target beo-http.service beo-server.service beo-sniffer.service

[Service]
Type=simple
User=kirsten
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/kirsten/.Xauthority
WorkingDirectory=/home/kirsten/beosound5c
ExecStart=/home/kirsten/beosound5c/webclient.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
